# nginx/files/dhparam-nginx.timer
[Unit]
Description=nginx dhparam regeneration timer

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
