# letsencrypt/files/letsencrypt.timer
[Unit]
Description=Letsencrypt certificate renewal timer

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
